# supervisor/supervisor.d/privoxy.ini
[program:privoxy]
command=/usr/local/sbin/privoxy --no-daemon /usr/local/etc/privoxy/config
autorestart=true
startretries=3
